# docker/supervisord.conf
; HTTP is served by the `app` service (Laravel Octane + Swoole).
; Horizon handles queue workers. The `scheduler` Compose service runs `schedule:work`.
; Start with: supervisord -n -c /etc/supervisor/supervisord.conf

[program:horizon]
process_name=%(program_name)s
command=php /var/www/html/artisan horizon
autostart=true
autorestart=true
stopwaitsecs=3600
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
